# mysql/crontab
0 2 * * * sh /backup.sh >> /var/log/backup.log 2>&1
